# rc.local
su - pi -c "THINGSPEAK_KEY='XXXXXXXXXXXXXXXX' ~/raspberry-teleinfo/start.sh"
