# rc.local
su - pi -c "THINGSPEAK_KEY='XXXXXXXXXXXXXXXX' CHANNEL_ID='XXXXX' ~/raspberry-teleinfo/start.sh"
